# systemd/ddns-cloudflare.service
[Unit]
Description=DDNS Updater

[Service]
ExecStart=/usr/bin/docker run --rm --name ddns-cloudflare-updater ddns-cloudflare:latest